# ethbuilder-operational/systemd/hermes-janitor.service
[Unit]
Description=Hermes tmp/chromium janitor
[Service]
Type=oneshot
ExecStart=/root/.hermes/scripts/hermes-janitor.sh
