# ethbuilder-operational/systemd/permfix-sentinel.service
[Unit]
Description=Permfix drift sentinel
[Service]
Type=oneshot
ExecStart=/root/.hermes/scripts/permfix-sentinel.sh
